# configs/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D


[program:jupyterhub]
command=jupyterhub  -f /opt/conda/jupyterhub_config.py 
                       
